# rootfs/etc/services.d/flexget-inotify/run
#!/usr/bin/with-contenv sh

cd $FLEXGET_SYNC_DIR
/usr/sbin/reflex -g '*.*' /usr/sbin/sync.sh